# hub/appliance/etc/systemd/user/fadecandy-watcher.path
[Path]
PathModified=/usr/local/bin

[Install]
WantedBy=multi-user.target