# appliance/etc/systemd/system/fadecandy-watcher.service
[Unit]
Description=FadeCandy restart config

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart fadecandy.service
