# appliance/etc/systemd/system/fadecandy-config.service
[Unit]
Description=Update FadeCandy config from universe memory

[Service]
ExecStart=/scripts/sync-config.sh
User=blinkenlights
