# appliance/etc/systemd/system/fadecandy-watcher.path
[Unit] 
Wants=fadecandy-watcher.service

[Path]
PathChanged=/usr/local/bin/config/fcserver.json

[Install]
WantedBy=multi-user.target
